# configs/osquery/osquery.conf
{
  "schedule": {
    "kernel_modules": {
      "query": "SELECT name, size, used_by, status FROM kernel_modules;",
      "interval": 300
    },
    "listening_ports": {
      "query": "SELECT pid, port, protocol, address FROM listening_ports;",
      "interval": 60
    },
    "process_events": {
      "query": "SELECT pid, path, cmdline, time FROM process_events;",
      "interval": 120
    },
    "suid_binaries": {
      "query": "SELECT path FROM suid_binaries;",
      "interval": 86400
    },
    "arp_cache": {
      "query": "SELECT address, mac, interface FROM arp_cache;",
      "interval": 300
    },
    "crontab": {
      "query": "SELECT * FROM crontab;",
      "interval": 3600
    },
    "authorized_keys": {
      "query": "SELECT * FROM authorized_keys;",
      "interval": 3600
    }
  },
  "file_paths": {
    "system_binaries": ["/usr/bin/%%", "/usr/sbin/%%"],
    "etc": ["/etc/%%"]
  }
}
